# crontab
*/2 * * * *   /root/monitoring.sh >> /root/monitoring.log
